# 90-hdmi.rules
# link to /etc/udev/rules.d
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/zach/.dotfiles/hdmi-switch.sh"
